# config/polis-init.service
[Unit]
Description=Polis Workspace Network Initialization
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/polis-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
