# agents/openclaw/config/openclaw.service
# =============================================================================
# OpenClaw Gateway - Systemd Service Unit
# =============================================================================
# Runs the OpenClaw AI agent gateway within the Polis workspace.
# Gateway provides Control UI and handles all agent interactions.
# =============================================================================

[Unit]
Description=OpenClaw AI Agent Gateway
Documentation=https://docs.openclaw.ai
After=network-online.target polis-init.service
Wants=network-online.target
Requires=polis-init.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=polis
Group=polis
WorkingDirectory=/app

# Load environment from init script output
EnvironmentFile=-/home/polis/.openclaw/.env

# Force IPv4 for DNS resolution (IPv6 is blocked in the container)
# Use full OS CA bundle for outbound TLS (OpenAI/API providers behind enterprise proxies)
Environment=NODE_OPTIONS=--dns-result-order=ipv4first
Environment=NODE_EXTRA_CA_CERTS=/etc/ssl/certs/ca-certificates.crt
Environment=PATH=/home/polis/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Inject Polis MITM CA into system trust store (resets on container restart)
# This ensures ALL TLS clients (undici, node-fetch, curl) trust the proxy CA
# The '+' prefix runs this step as root (needed to write to /etc/ssl/certs/)
# We check by counting certs — if polis CA is missing, append it
ExecStartPre=+/bin/bash -c 'if ! openssl verify -CAfile /etc/ssl/certs/ca-certificates.crt /etc/ssl/certs/polis-ca.pem >/dev/null 2>&1; then cat /etc/ssl/certs/polis-ca.pem >> /etc/ssl/certs/ca-certificates.crt && echo "[openclaw] Polis CA injected into system trust store"; else echo "[openclaw] Polis CA already trusted"; fi'

# Run initialization before starting gateway
ExecStartPre=/usr/local/bin/openclaw-init.sh

# Wait for transparent proxy to be fully ready (TPROXY routes + ICAP)
# Use -k to skip TLS verify for the check itself, and -o /dev/null -w '%{http_code}'
# to accept any non-connection-error response (even 401) as "proxy is working"
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do code=$(curl -s -o /dev/null -w "%%{http_code}" --max-time 5 https://api.openai.com/v1/models 2>/dev/null); [ "$code" -gt 0 ] 2>/dev/null && exit 0; sleep 2; done; echo "WARNING: proxy not ready after 60s, starting anyway"'

# Start OpenClaw gateway
# --allow-unconfigured: Start even without full onboarding
# --bind lan: Bind to 0.0.0.0 (accessible within container)
# --port 18789: Gateway port for Control UI
ExecStart=/usr/bin/node dist/index.js gateway --allow-unconfigured --bind lan --port 18789

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=10

# Security hardening (relaxed for container compatibility)
# NoNewPrivileges, ProtectSystem, ProtectHome disabled — container lacks
# required capabilities for these sandboxing features.
ReadWritePaths=-/home/polis/.openclaw

# Resource limits
MemoryMax=2G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
